# scripts/com.darvin.agentic-loops-queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.darvin.agentic-loops-queue</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/darvin/Documents/Projects/darvinyi-agentic-loops/scripts/queue-worker.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/darvin</string>
    <key>PATH</key><string>/Users/darvin/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>StartInterval</key><integer>60</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/private/tmp/agentic-loops-queue.stdout.log</string>
  <key>StandardErrorPath</key><string>/private/tmp/agentic-loops-queue.stderr.log</string>
</dict>
</plist>
