# scripts/com.darvin.agentic-loops-queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.darvin.agentic-loops-queue</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/darvin/Documents/Projects/darvinyi-agentic-loops/scripts/queue-worker.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key><string>/Users/darvin</string>
    <key>PATH</key><string>/Users/darvin/Documents/Projects/darvinyi-agentic-loops/scripts/service-bin:/Users/darvin/.local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/usr/local/bin</string>
    <key>PIPELINE_GIT_BIN</key><string>/usr/bin/git</string>
    <key>PIPELINE_GIT_TIMEOUT_SECONDS</key><string>60</string>
    <key>TERRA_IDLE_TIMEOUT_SECONDS</key><string>1800</string>
    <key>TERRA_MAX_RUNTIME_SECONDS</key><string>5400</string>
    <key>PIPELINE_VALIDATION_IDLE_TIMEOUT_SECONDS</key><string>600</string>
    <key>PIPELINE_VALIDATION_MAX_RUNTIME_SECONDS</key><string>1800</string>
    <key>PIPELINE_SYNC_IDLE_TIMEOUT_SECONDS</key><string>120</string>
    <key>PIPELINE_SYNC_MAX_RUNTIME_SECONDS</key><string>900</string>
    <key>PIPELINE_WATCHDOG_TERM_GRACE_SECONDS</key><string>10</string>
    <key>PIPELINE_WATCHDOG_POLL_SECONDS</key><string>2</string>
  </dict>
  <key>StartInterval</key><integer>60</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/private/tmp/agentic-loops-queue.stdout.log</string>
  <key>StandardErrorPath</key><string>/private/tmp/agentic-loops-queue.stderr.log</string>
</dict>
</plist>
